# root/etc/services.d/cron/run
#!/usr/bin/with-contenv sh

if [[ "$CRON" ]]; then
  crond -l 2 -f
else
  echo "Skipping cron, all done."
  exit 0
fi

