# crontab
0 4 * * * /run.sh "$(cat /bucket_name.txt)"
